# netbox-scaffold/netbox/templates/netbox-sync.timer
[Unit]
Description=Hourly vCenter refresh

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
